# etc/rc.local
#!/bin/sh -e

apachectl -k start

exit 0
